# etc/cron.daily/vmbackup
#!/bin/sh

test -x /usr/local/sbin/vmbackuprotate || exit 0
/usr/local/sbin/vmbackuprotate /etc/vmbackup.conf >/dev/null
